# scripts/launchd/com.randomify.prod-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- prod weekly refresh: full (uncapped) backlog rebuild from the latest dump. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.randomify.prod-refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/seanreid/randomify/scripts/refresh-dump.sh</string>
    <string>production</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/seanreid/randomify/data/musicbrainz/logs/prod-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/seanreid/randomify/data/musicbrainz/logs/prod-refresh.log</string>
</dict>
</plist>
